# iprediaos-builder.cron
#! /bin/bash

curl https://raw.github.com/mattias-ohlsson/iprediaos-builder/master/iprediaos-build.sh | bash